# tools/com.fpl.transfer-poll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fpl.transfer-poll</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/ruimiguelneves/Code/fpl-dashboard/tools/poll_transfers.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>WorkingDirectory</key>
    <string>/Users/ruimiguelneves/Code/fpl-dashboard</string>

    <key>StandardOutPath</key>
    <string>/Users/ruimiguelneves/Code/fpl-dashboard/tools/cache/poll.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ruimiguelneves/Code/fpl-dashboard/tools/cache/poll.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/ruimiguelneves/Library/Python/3.9/bin</string>
        <key>HOME</key>
        <string>/Users/ruimiguelneves</string>
    </dict>
</dict>
</plist>
